# dev-servers/dev-desktop-x11vnc/clients/supervisord.conf
[supervisord]
nodaemon=true

[program:novnc]
command=/root/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autorestart=true

[program:webssh]
command=wssh
autorestart=true